# etc/supervisord.d/inventory.ini
[program:inventory]
command=/opt/qnib/inventory/bin/inventory.py --server -l
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
